# root/etc/uci-defaults/luci-homeproxy-migration
#!/bin/sh

exec ucode "/etc/homeproxy/scripts/migrate_config.uc"
